# backend/com.jandrishti.aqi.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.jandrishti.aqi</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Library/Frameworks/Python.framework/Versions/3.13/bin/python3</string>
        <string>/Users/prativa/Documents/GitHub/JanDrishti-Pollution-Action-DashBoard/backend/run_server.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/prativa/Documents/GitHub/JanDrishti-Pollution-Action-DashBoard/backend</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/prativa/Documents/GitHub/JanDrishti-Pollution-Action-DashBoard/backend/logs/aqi_service.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/prativa/Documents/GitHub/JanDrishti-Pollution-Action-DashBoard/backend/logs/aqi_service.error.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
